# updater.cron
MAILTO=michal.vaner@nic.cz
0     *       *       *       *       /usr/bin/updater.sh
