# updater.cron
MAILTO=michal.vaner@nic.cz
0     *       *       *       *       root	/usr/bin/updater.sh
