# tests/system/install/output/etc/cron.d/updater
MAILTO=""
0	*       *       *       *       root	/usr/bin/updater.sh >/dev/null 2>&1 &
0	0	*	*	*	root	/usr/bin/updater-unstuck.sh >/dev/null 2>&1 &
