# tests/system/steal-confs.sys_update/output/usr/lib/opkg/info/dropbear.list
/etc/config/dropbear
/etc/dropbear/dropbear_dss_host_key
/etc/dropbear/dropbear_rsa_host_key
/etc/init.d/dropbear
/usr/bin/dbclient
/usr/bin/dropbearkey
/usr/bin/scp
/usr/bin/ssh
/usr/sbin/dropbear
